# backend/phidro.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd LaunchDaemon — Pedal Hidrografico (servir o backend localmente
  no macOS). TROQUE: o usuario e os caminhos.
  Instalar: ver a secao "Rodar no macOS" do README.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>co.pedalhidrografi.phidro</string>

  <key>UserName</key>
  <string>danilo</string>

  <key>WorkingDirectory</key>
  <string>/Users/danilo/pedalhidrografico/backend</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PHIDRO_DATA</key>
    <string>/Users/danilo/phidro-data</string>
  </dict>

  <!-- workers=1 OBRIGATORIO: o lock global de mutacao e per-processo;
       2+ workers = lost updates nos catalogos TTL. -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/danilo/pedalhidrografico/backend/venv/bin/gunicorn</string>
    <string>--bind</string>
    <string>127.0.0.1:8000</string>
    <string>--workers</string>
    <string>1</string>
    <string>--threads</string>
    <string>8</string>
    <string>--timeout</string>
    <string>120</string>
    <string>main:app</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/phidro.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/phidro.err</string>
</dict>
</plist>
