# backend/pi/phidro.service
[Unit]
Description=Pedal Hidrografico — backend e app do mapa
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pedalhidrografico/backend/pi
Environment=PHIDRO_DATA=/home/pi/phidro-data
# --workers 1 OBRIGATORIO: o lock global de mutacao (_state_lock) e
# per-processo; com 2+ workers, uploads concorrentes (o form manda cards
# em paralelo) caem em processos diferentes e o segundo write sobrescreve
# os triples do primeiro (lost update). Threads cobrem a concorrencia.
ExecStart=/home/pi/pedalhidrografico/backend/pi/venv/bin/gunicorn \
  --bind 127.0.0.1:8000 --workers 1 --threads 8 --timeout 120 main:app
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
